# phps/8.0/app.conf
[program:app_queue]
command=php /var/www/html/artisan queue:work --sleep=3 --tries=3 --queue=high,medium,low,default
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
redirect_stderr=true
